# crontab
0 * * * * automate /usr/local/bin/python /app/script.py >> /var/log/cron.log 2>&1